# trunk/se3-logonpy/etc/cron.d/se3-logonpy
# Rajout  disable-gpo.lck pour les postes faisant partis d'un parc ayant un parc avec "description: disable_gpo"
0 6  * * *   root    /usr/share/se3/sbin/no-gpo-upload.sh >/dev/null 2>&1
